# roda-ui/roda-wui/docker-files/supervisor-conf.d/clamav.conf
[program:clamav]
command=/usr/sbin/clamd -c /etc/clamav/clamd.conf
