# roda-ui/roda-wui/docker-files/supervisor-conf.d/siegfried.conf
[program:siegfried]
command=sf -serve localhost:5138
